# rivian-gearshop-crawler.service
[Unit]
Description=Rivian Gear Shop Crawler
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/rivian-gearshop-crawler
ExecStart=/opt/rivian-gearshop-crawler/venv/bin/python3 /opt/rivian-gearshop-crawler/crawler.py
Environment=PATH=/opt/rivian-gearshop-crawler/venv/bin:/usr/bin:/bin
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal
